# crontab
*/1 * * * * php /slackcamp/slackcamp.php
